# units/jenkins@.service.d/alpine.conf
# Use the alpine base image instead of the default one.
[Service]
Environment='DOCKER_IMAGE=jenkins:alpine'
