# units/jenkins@.service.d/publish-50000.conf
[Service]
# Expose the slave agent port (50000)
Environment='DOCKER_PUBLISH_50000=50000:50000'
