# units/jenkins@.service.d/virtual-port-8080.conf
# Sets the VIRTUAL_PORT environment variable to 8080.
#
# This provides compatability with the nginx-proxy.
[Service]
Environment='DOCKER_ENV_VIRTUAL_PORT=VIRTUAL_PORT=8080'
